# tools/pve-install-agent/sideloop-install-agent.timer
# /etc/systemd/system/sideloop-install-agent.timer
# Vérifie le manifest toutes les 30 min. L'agent est idempotent PAR DEVICE
# (installed-state.json : un device servi il y a < 48 h est skippé) → il n'installe
# que les devices en retard, et seulement quand l'iPhone est joignable.
[Unit]
Description=sideloop install agent (toutes les 30 min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
